# rootfs/etc/s6-overlay/s6-rc.d/cups-server/run
#!/command/with-contenv bash
ulimit -n 1048576 2>/dev/null || true

# Wait for host Avahi to be reachable via D-Bus (host_dbus: true mounts host socket)
until dbus-send --system --dest=org.freedesktop.Avahi --print-reply \
    /org/freedesktop/Avahi org.freedesktop.Avahi.Server.GetVersionString \
    >/dev/null 2>&1; do
    sleep 1
done

exec cupsd -f
