# rootfs/etc/s6-overlay/s6-rc.d/avahi-daemon/run
#!/command/with-contenv bash
until [ -e /run/dbus/system_bus_socket ]; do
    sleep 1
done
exec avahi-daemon
