# rootfs/etc/s6-overlay/s6-rc.d/avahi-daemon/run
#!/command/with-contenv bash
exec avahi-daemon --no-drop-root --no-chroot
