# rootfs/etc/s6-overlay/s6-rc.d/avahi-daemon/run
#!/command/with-contenv bash
# Wait for D-Bus socket to be ready (avahi needs it)
until [ -e /var/run/dbus/system_bus_socket ]; do
    sleep 1
done
exec avahi-daemon --no-drop-root --no-chroot
