# rootfs/etc/s6-overlay/s6-rc.d/dbus-daemon/run
#!/command/with-contenv bash
rm -f /var/run/dbus/pid
exec dbus-daemon --system --nofork
